# deploy/systemd/sohotexto.service
[Unit]
Description=Justext — content from youtube
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sohotexto
ExecStart=/usr/bin/node server.mjs
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=SOHOTEXTO_HOST=127.0.0.1
Environment=SOHOTEXTO_PORT=3217
# ANTHROPIC_API_KEY must be set here — never commit the real value
Environment=ANTHROPIC_API_KEY=sk-ant-...

[Install]
WantedBy=multi-user.target
